# mkosi.extra/usr/lib/systemd/system/systemd-sysupdate.service.d/notify-pending.conf
[Unit]
OnSuccess=system-update-pending.service
