# mkosi.extra/usr/lib/systemd/system/system-update-pending.service
[Unit]
Description=Notify user of pending system update

[Service]
Type=oneshot
ExecStart=/usr/bin/system-update-pending

[Install]
WantedBy=graphical.target
